# startup/crontab/crontab.entry
@reboot sleep 60 && ~/crontab/start_procServ_EPICS_apfel_IOC.sh
