# startup/crontab/start_procServ_EPICS_apfel_IOC.sh
#!/bin/bash
#procServ -q -L /home/panda/procServ/start_apfel_IOC_procServ.log --logstamp -n apfelControlIOC -l 4814 4813 /home/panda/procServ/start_apfel_IOC_procServ.sh 
procServ -q -n apfelControlIOC -l 4814 4813 /home/panda/procServ/start_apfel_IOC_procServ.sh 

